# cli/systemd/safescale.socket
[Socket]
ListenStream=127.0.0.1:50051

[Install]
WantedBy=sockets.target
